# pkg/network/templates/setup-3fs-network.service.tmpl
[Unit]
Description=Setup 3fs network
Before=docker.service
 
[Service]
Type=simple
ExecStart={{.ScriptPath}}
 
[Install]
WantedBy=multi-user.target