# pkg/storage/templates/mount_3fs_disks.service.tmpl
[Unit]
Description=Mount 3fs disks
Before=docker.service
 
[Service]
Type=simple
ExecStart={{.ScriptPath}} {{.BaseDir}}
 
[Install]
WantedBy=multi-user.target